# launchd/com.driftdriver.zombie-reaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.driftdriver.zombie-reaper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/driftdriver</string>
    <string>reaper</string>
    <string>run</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/driftdriver-zombie-reaper.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/driftdriver-zombie-reaper.err</string>
</dict>
</plist>
